# express/supervisor.d/consul-template.ini
[program:consul-template]
command=/bin/consul-template -config=/etc/consul-template.conf
autostart=false
autorestart=false
stdout_logfile=/tmp/%(program_name)s.log
stderr_logfile=/tmp/%(program_name)s.log
