# vault/supervisor.d/consul-template.ini
[program:consul]
command=/bin/consul-template -config=/etc/consul-template.d/config.hcl
autostart=true
autorestart=false
stdout_logfile=/tmp/%(program_name)s.log
stderr_logfile=/tmp/%(program_name)s.log
