# vault/supervisor.d/consul.ini
[program:consul]
command=/bin/consul agent -config-dir=/etc/consul.d -dns-port=53
autostart=true
autorestart=false
stdout_logfile=/tmp/%(program_name)s.log
stderr_logfile=/tmp/%(program_name)s.log
